# initramfs/basic/etc/inittab
::sysinit:/etc/init.d/rcS
console::respawn:/bin/sh
::restart:/sbin/init
