# initramfs/basic/etc/inittab
::sysinit:/etc/init.d/rcS
console::respawn:-/bin/sh
#ttyUL0::respawn:/sbin/getty -L ttyUL0 115200 vt100
::restart:/sbin/init

